# foss-setup/configs/ansible/ansible-pull.service
# ansible-pull.service — self-converging fleet, the PULL half of glue-08.
#
# Each host clones the control repo and runs site.yml against ITSELF on a timer,
# so a drifted or rebuilt box returns to its defined state with nobody driving
# it. Install on every Ansible-managed host (Mac mini, seedbox, rig):
#
#   sudo cp ansible-pull.service ansible-pull.timer /etc/systemd/system/
#   sudo systemctl enable --now ansible-pull.timer
#
# Edit CONTROL_REPO_URL below to your Forgejo control-repo SSH URL (glue-05/06).
[Unit]
Description=ansible-pull fleet convergence (glue-08)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=btabaska
Group=btabaska
Environment=PATH=/home/btabaska/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=CONTROL_REPO_URL=git@forgejo:home/homelab.git
Environment=ANSIBLE_CONFIG=/home/btabaska/.ansible-pull/configs/ansible/ansible.cfg

ExecStart=/bin/bash -c '/home/btabaska/.local/bin/ansible-pull \
    -U "${CONTROL_REPO_URL}" \
    -d /home/btabaska/.ansible-pull \
    --accept-host-key \
    -i configs/ansible/inventory.ini \
    --limit "$(hostname -s)" \
    --connection local \
    --diff --check \
    configs/ansible/site.yml'

# Don't let a transient pull failure mark the box permanently failed.
SuccessExitStatus=0

[Install]
WantedBy=multi-user.target
