# foss-setup/configs/ansible/ansible-pull.timer
# ansible-pull.timer — schedule for the convergence run (glue-08).
#
# Always-on hosts (Mac mini, seedbox): converge a few times a day, with jitter
# so the fleet doesn't all hammer Forgejo at once.
#
# On-demand rig (cachyos): it sleeps through any wall-clock schedule, so DROP
# OnCalendar there and use the wake hook instead — set Persistent=true (run the
# missed job on next wake) AND add an OnBootSec so it converges shortly after
# Wake-on-LAN (game-08), matching wake_gated in group_vars/on_demand.yml.
[Unit]
Description=Run ansible-pull fleet convergence on a schedule (glue-08)

[Timer]
OnCalendar=*-*-* 04:20:00
RandomizedDelaySec=30m
Persistent=true
# Rig only: uncomment to converge ~3 min after each wake instead of by clock.
# OnBootSec=3min

[Install]
WantedBy=timers.target
