# foss-setup/configs/ansible/ansible-pull.timer
# ansible-pull.timer — schedule for the convergence run (glue-08).
#
# All hosts (Mac mini, seedbox, rig): converge a few times a day, with jitter
# so the fleet doesn't all hammer Forgejo at once. The rig runs 24/7 as of
# 2026-07-08, so it uses the same OnCalendar schedule as everyone else.
# Persistent=true is generic catch-up (run a missed job after downtime) — it is
# not a wake-gating mechanism.
[Unit]
Description=Run ansible-pull fleet convergence on a schedule (glue-08)

[Timer]
OnCalendar=*-*-* 04:20:00
RandomizedDelaySec=30m
Persistent=true
# Optional: converge ~3 min after boot too — useful after a recovery boot
# (power outage / accidental shutdown), harmless otherwise.
# OnBootSec=3min

[Install]
WantedBy=timers.target
