# foss-setup/configs/host/rig/music-mirror/music-mirror.service
[Unit]
Description=Mirror NAS music library to ~/Music (rsync --delete; Rhythmbox/iPod source)
# ~/Music is a MIRROR of the NAS master — additions AND deletions propagate,
# so local duplicates/drift are structurally impossible (2026-07-10: the old
# hand-copied ~/Music had months of drift plus the library's orphan dupes).
# NAS master is the only place music is ever added (via Lidarr).
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=btabaska
# ro automount; trailing slashes = mirror contents. CIFS has no owners/links,
# so -rt (not -a). --delete-after keeps the window of missing files minimal.
ExecStart=/usr/bin/rsync -rt --delete-after --exclude '#recycle' --exclude '@eaDir' /mnt/nas-music-ro/ /home/btabaska/Music/
# Dead-man: healthchecks "music-mirror-rig". Deploy note: replace placeholder
# with vault healthchecks.music_mirror_rig_ping_url
ExecStartPost=/usr/bin/curl -fsS -m 10 --retry 3 -o /dev/null HC_MUSIC_MIRROR_PING_URL
TimeoutStartSec=30min
Nice=15
