# foss-setup/configs/host/rig/music-mirror/nas-music-mirror.service
[Unit]
Description=Mirror NAS music library -> ~/Music as iPod-playable ALAC/MP3
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=btabaska
Nice=10
IOSchedulingClass=idle
ExecStart=/home/btabaska/bin/nas-music-to-alac-mirror.sh
