# foss-setup/configs/host/rig/ai-stack-watchdog/ai-stack-watchdog.service
[Unit]
Description=AI stack watchdog: container->host ollama hop -> healthchecks dead-man
After=docker.service ollama.service
Wants=docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/ai-stack-watchdog.env
ExecStart=/usr/local/bin/ai-stack-watchdog.sh
