# foss-setup/configs/host/rig/pcie-aer-monitor/pcie-aer-monitor.timer
[Unit]
Description=Run PCIe AER monitor every 20 min

[Timer]
OnBootSec=5min
OnUnitActiveSec=20min
Persistent=true

[Install]
WantedBy=timers.target
