# foss-setup/scripts/backup/restic-backup.timer
# Fires the restic backup nightly at 01:30 on every host — the rig included
# (24/7 since 2026-07-08). Persistent=true is generic catch-up: it runs a
# missed job at next boot if a host happened to be down at the scheduled time.
[Unit]
Description=Run restic backup daily

[Timer]
OnCalendar=*-*-* 01:30:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
