# foss-setup/verification/systemd/verification-quick.service
[Unit]
Description=Hourly quick verification tier (url checks: rig AI e2e, edge vhosts, DNS)
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
User=btabaska
EnvironmentFile=-/etc/verification/env
# --respect-enabled: this is a SCHEDULED tier, so honor `enabled` — never
# resurrect a deliberately-disabled check the way an ad-hoc operator `--host`
# run does (M-quirk: the quick tier would silently run + page disabled checks).
# --host url = the cheap externally-probeable tier (includes rig-ai-e2e, the
# real-completion check added after the 2026-07-09 UFW "no models" regression
# sat invisible for hours because the full sweep only runs daily at 07:15).
# --notify: quick tier alerts on state change; ad-hoc filtered runs stay silent.
ExecStart=/opt/verification/bin/run-checks.sh --host url --notify --respect-enabled
# Fleet-wide container manifests + health + failed systemd units (see
# checks.d/docker-fleet.yaml) — the blanket "anything stopped working" layer.
ExecStart=/opt/verification/bin/run-checks.sh --host docker-fleet --notify --respect-enabled
# Media library data-quality (checks.d/media.yaml): duplicate-album tripwire.
ExecStart=/opt/verification/bin/run-checks.sh --host media --notify --respect-enabled
# SuccessExitStatus=1: a crit check failing exits run-checks.sh 1, but the tier
# still RAN — treat exit 1 as success so (a) the next ExecStart still runs, (b)
# the unit isn't marked FAILED (which the systemd-failed-mini check would then
# flag against the verification framework itself — a self-referential alert),
# and (c) the dead-man ping below still fires. A real crash/timeout exits 2+ and
# correctly fails the unit + downs the dead-man. (Same fix as verification-fast.)
SuccessExitStatus=1
# Dead-man: healthchecks "verification-quick-mini" (ntfy on silence >1h30m).
# URL is in /etc/verification/env as VERIFY_QUICK_PING_URL (kept out of git so
# the repo unit installs verbatim — no hand-substituted placeholder).
ExecStartPost=/usr/bin/curl -fsS -m 10 --retry 3 -o /dev/null ${VERIFY_QUICK_PING_URL}
TimeoutStartSec=15min
Nice=10

[Install]
WantedBy=multi-user.target
