# foss-setup/configs/host/mini/net-selfheal/net-selfheal.timer
[Unit]
Description=Run net-selfheal every minute (recover NIC after DHCP lease loss)

[Timer]
OnBootSec=90
OnUnitActiveSec=60
AccuracySec=10s

[Install]
WantedBy=timers.target
