# foss-setup/scripts/inventory/export-manifests.timer
# export-manifests.timer — run export-manifests.service weekly.
# Install alongside export-manifests.service in /etc/systemd/system/.
#
#   sudo systemctl enable --now export-manifests.timer
#   systemctl list-timers export-manifests.timer

[Unit]
Description=Weekly host manifest export into the control repo
Documentation=https://etckeeper.branchable.com/

[Timer]
# Weekly: Mondays at 04:00 local time.
OnCalendar=Mon *-*-* 04:00:00
# Catch up if the box was off at the scheduled time.
Persistent=true
RandomizedDelaySec=20m

[Install]
WantedBy=timers.target
