# foss-setup/scripts/inventory/sbom-nightly.service
# sbom-nightly.service — generate + upload SBOMs to Dependency-Track (oneshot).
# Triggered by sbom-nightly.timer. Install to /etc/systemd/system/.
#
#   sudo cp sbom-nightly.service sbom-nightly.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now sbom-nightly.timer
#
# Secrets (DTRACK_URL, DTRACK_API_KEY) come from the EnvironmentFile below —
# create it from sbom.env.example and chmod 600 it.

[Unit]
Description=Nightly SBOM generation + upload to Dependency-Track
Documentation=https://docs.dependencytrack.org/usage/cicd/
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
Environment=PATH=/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/etc/inventory/sbom.env
ExecStart=/opt/scripts/inventory/sbom-nightly.sh
# Scanning the whole filesystem + pulling image layers is heavy; be patient/nice.
TimeoutStartSec=2h
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
